# services/smf-manifest/tomcat-jenkins.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!-- NOTE: Jenkins requires either Oracle JDK or maybe OpenJDK+fontconfig+TTF fonts -->
<service_bundle type='manifest' name='export'>
  <service name='application/tomcat' type='service' version='0'>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    <instance name='jenkins' enabled='false'>
      <method_context project=':default' resource_pool=':default' working_directory=':default'>
        <method_credential limit_privileges=':default'
         privileges='basic,net_privaddr' supp_groups=':default'
         group=':default' user='root'/>
        <!-- NOTE: this privilege set allows an unprivileged user account
             to utilize privileged ports like HTTP = 80/tcp -->
      </method_context>
      <exec_method name='restart' type='method' exec='/etc/init.d/cataliner.sh %m' timeout_seconds='0'>
        <method_context/>
      </exec_method>
      <exec_method name='start' type='method' exec='/etc/init.d/cataliner.sh %m' timeout_seconds='0'>
        <method_context/>
      </exec_method>
      <exec_method name='stop' type='method' exec='/etc/init.d/cataliner.sh %m' timeout_seconds='0'>
        <method_context/>
      </exec_method>
      <property_group name='cataliner' type='application'>
        <propval name='CATALINER_APPSERVER_DIR' type='astring' value='/opt/tomcat'/>
        <!-- propval name='JAVA_HOME' type='astring' value='/opt/java'/ -->
        <propval name='CATALINER_SMF_CONFIGFILE_EXCLUSIVE' type='boolean' value='true'/>
        <propval name='CATALINER_SMF_CONFIGFILE_PATH' type='astring' value='/etc/default/jenkins-tomcat.conf'/>
      </property_group>
    </instance>
    <template>
      <common_name>
        <loctext xml:lang='C'>Apache Tomcat, Java Servlet Container, which hosts Magnolia website content-management system</loctext>
      </common_name>
      <documentation>
        <doc_link name='apache.org' uri='http://tomcat.apache.org/'/>
      </documentation>
    </template>
  </service>
</service_bundle>
